# systemd/network-scan.service
# Network Scan Service (Alternative to cron)
#
# This systemd service runs the network scanner periodically.
# Use this instead of cron if you prefer systemd timers.
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp network-scan.service ~/.config/systemd/user/
#   cp network-scan.timer ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now network-scan.timer
#
# Check status:
#   systemctl --user status network-scan.timer
#   systemctl --user list-timers
#

[Unit]
Description=Network Scanner Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/bin/network-scan-daemon
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
